# deploy/rubika-panel.service
[Unit]
Description=Rubika Userbot Management Panel (Telegram-controlled)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Change User to your server user (e.g. ubuntu / root)
User=root
# Absolute path to the project folder on the server
WorkingDirectory=/opt/rubika-panel
# Use the venv python so dependencies are found
ExecStart=/opt/rubika-panel/venv/bin/python /opt/rubika-panel/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
